# crontab
#run python script every minute
SHELL=/bin/bash
BASH_ENV=/container.env
10 11 * * * python3 /app/xml_to_turtle.py > /proc/1/fd/1 2>/proc/1/fd/2
